# deb/DEBIAN/postrm
#!/bin/sh

set -e

update-rc.d -f hsp-fo-app remove
